# crontab
@reboot python /home/pi/meteo/src/add_to_files.py &
